# build/frt.service
[Unit]
Description=frt

[Service]
ExecStart=/usr/bin/make -C /build

[Install]
WantedBy=multi-user.target